# etc/linux/wa-mini-update.service
[Unit]
Description=WhatsApp version update check

[Service]
Type=oneshot
User=_wamini
Group=_wamini
ExecStart=/usr/local/bin/wa-mini update-version
WorkingDirectory=/var/lib/wa-mini

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ReadWritePaths=/var/lib/wa-mini
